# psyrcd.service
[Unit]
Description=Psyrcd The Psybernetics IRC Server

[Service]
ExecStart=/usr/bin/psyrcd

[Install]
WantedBy=multi-user.target
